# openwrt-sdk-ar71xx-generic_gcc-7.3.0_musl.Linux-x86_64/staging_dir/target-mips_24kc_musl/root-ar71xx/etc/hotplug.d/usbmisc/01_get_signal.sh
#!/bin/sh

[ "$ACTION" = add ] || [ "$ACTION" = remove ] || exit 0
[ "${DEVNAME/[0-9]/}" = cdc-wdm ] || exit 0

[ "$ACTION" = add ] && modem_get_signal /dev/$DEVNAME  &
[ "$ACTION" = remove ] && { 
	rm /tmp/signal
	pid=$(fuser  /usr/bin/modem_get_signal)
	kill $pid 
	}

